# infra/systemd/rapidrag.service
# /etc/systemd/system/rapidrag.service
# RapidRAG Docker Compose Stack — Phased Startup
#
# WHY THIS EXISTS:
# Ensures the full compose stack (Vault, n8n, PostgreSQL, Redis, MinIO, etc.)
# starts automatically after Docker is fully ready on every server boot.
# Uses a phased startup script to start containers in dependency order,
# preventing health check failures and disk I/O spikes from simultaneous startup.
#
# DESIGN DECISION — No hard Requires=docker.service:
# Using Requires= causes systemd to fail this service immediately if Docker is
# still "activating" (which happens for 5-10 minutes on this HDD due to BoltDB
# fdatasync). Instead we use After= (ordering only) and poll for Docker readiness
# inside ExecStartPre, allowing indefinite wait with automatic retry.
#
# STARTUP PHASES (see infra/scripts/phased-startup.sh for details):
#   Phase 1: vault
#   Phase 2: vault-bootstrap + all *-vault-agent sidecars + cert-rotation-controller
#   Phase 3: postgres, redis, rabbitmq, minio, keycloak, opensearch, n8n, dify-db, dify-redis, dify-sandbox, dify-migrate
#   Phase 4: db-migrate → dify-api + dify-worker
#   Phase 5: workflow-service, logging-service, dify-web
#   Phase 6: api-gateway → web → web-ingress
#
# INSTALL (one-time on production server):
#   sudo cp /home/bali/09_rapidrag/infra/systemd/rapidrag.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable rapidrag   # enable auto-start on boot
#   sudo systemctl start rapidrag    # start now
#
# USAGE:
#   systemctl start rapidrag         # start the stack
#   systemctl stop rapidrag          # gracefully stop (preserves volumes)
#   systemctl status rapidrag        # check status
#   journalctl -u rapidrag -f        # follow startup progress

[Unit]
Description=RapidRAG Docker Compose Stack (Phased Startup)
Documentation=file:///home/bali/09_rapidrag/README.md
# After= ensures ordering (starts after docker if docker is running) but does NOT
# cause this service to fail if docker isn't active yet. We handle docker readiness
# ourselves in ExecStartPre with a polling loop.
After=docker.service containerd.service network-online.target
Wants=network-online.target docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/bali/09_rapidrag
Environment=ENVIRONMENT=dev

# Wait for Docker socket to be truly ready before starting.
# This loop polls up to 20 minutes (300 × 4s) — handles the BoltDB fdatasync
# delay on spinning HDD which takes 5-10 minutes. No hard dependency needed.
ExecStartPre=/bin/bash -c '\
  echo "Waiting for Docker to be ready..."; \
  for i in $(seq 1 300); do \
    docker info >/dev/null 2>&1 && \
    echo "Docker ready after $((i*4))s" && break || \
    sleep 4; \
  done; \
  docker info >/dev/null 2>&1 || { echo "Docker not ready after 20 minutes"; exit 1; }'

# Run the phased startup script via bash — avoids SELinux exec restrictions on /home
ExecStart=/bin/bash /home/bali/09_rapidrag/infra/scripts/phased-startup.sh

# On stop: bring down the stack gracefully (preserving volumes).
# Vault secrets are not needed for stop — placeholder env is sufficient.
ExecStop=/bin/bash -c 'docker compose \
  -f /home/bali/09_rapidrag/docker-compose.yml \
  -f /home/bali/09_rapidrag/docker-compose.prod.yml \
  --env-file /home/bali/09_rapidrag/.env.production \
  down --timeout 30 || true'

# Restart on failure with generous intervals
Restart=on-failure
RestartSec=120s

# Allow up to 30 minutes for all phases (Docker wait + dify-migrate ~12min + rest)
TimeoutStartSec=1800
TimeoutStopSec=300

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
